# env/auto-upgrade.timer
[Unit]
Description=Weekly FincLab env rebuild

[Timer]
OnCalendar=Sun 04:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
